# roles/node-exporter/templates/node_exporter.service.j2
[Unit]
Description=Prometheus Node Exporter
Documentation=https://github.com/prometheus/node_exporter
After=network.target

[Service]
User={{ node_exporter_user }}
Group={{ node_exporter_user }}
Type=simple
ExecStart={{ node_exporter_binary }}

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target